# src/var/lib/dpkg/info/ferm.list
/.
/etc
/etc/ferm
/etc/ferm/ferm.conf
/etc/ferm/ferm.d
/etc/init.d
/etc/init.d/ferm
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ferm.service
/usr
/usr/sbin
/usr/sbin/ferm
/usr/sbin/import-ferm
/usr/share
/usr/share/doc
/usr/share/doc/ferm
/usr/share/doc/ferm/AUTHORS
/usr/share/doc/ferm/NEWS.gz
/usr/share/doc/ferm/README.Debian
/usr/share/doc/ferm/README.rst
/usr/share/doc/ferm/TODO
/usr/share/doc/ferm/changelog.Debian.gz
/usr/share/doc/ferm/copyright
/usr/share/doc/ferm/examples
/usr/share/doc/ferm/examples/antiddos.ferm
/usr/share/doc/ferm/examples/arptables.ferm
/usr/share/doc/ferm/examples/dmz_router.ferm
/usr/share/doc/ferm/examples/dsl_router.ferm
/usr/share/doc/ferm/examples/ebtables.ferm
/usr/share/doc/ferm/examples/fileserver.ferm
/usr/share/doc/ferm/examples/ipv6.ferm
/usr/share/doc/ferm/examples/mailserver.ferm
/usr/share/doc/ferm/examples/resolve.ferm
/usr/share/doc/ferm/examples/webserver.ferm
/usr/share/doc/ferm/examples/workstation.ferm
/usr/share/doc/ferm/ferm.html
/usr/share/doc/ferm/ferm.txt.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/ferm
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/ferm.1.gz
/usr/share/man/man1/import-ferm.1.gz
/var
/var/cache
/var/cache/ferm
