# src/var/lib/dpkg/info/gcr.postrm
#!/bin/sh
set -e
# Automatically added by dh_installsystemduser/13.10
if [ "$1" = "purge" ]; then
	if [ -z "${DPKG_ROOT:-}" ] && [ -x "/usr/bin/deb-systemd-helper" ] ; then
		deb-systemd-helper --user purge 'gcr-ssh-agent.service' 'gcr-ssh-agent.socket' >/dev/null || true
	fi
fi
# End automatically added section
