# src/var/lib/dpkg/info/gnustep-base-runtime.list
/.
/etc
/etc/default
/etc/default/gdomap
/etc/init.d
/etc/init.d/gdomap
/usr
/usr/bin
/usr/bin/HTMLLinker
/usr/bin/autogsdoc
/usr/bin/cvtenc
/usr/bin/defaults
/usr/bin/gdnc
/usr/bin/gdomap
/usr/bin/gspath
/usr/bin/make_strings
/usr/bin/pl
/usr/bin/pl2link
/usr/bin/pldes
/usr/bin/plget
/usr/bin/plmerge
/usr/bin/plparse
/usr/bin/plser
/usr/bin/plutil
/usr/bin/sfparse
/usr/bin/xmlparse
/usr/share
/usr/share/doc
/usr/share/doc/gnustep-base-runtime
/usr/share/doc/gnustep-base-runtime/NEWS.Debian.gz
/usr/share/doc/gnustep-base-runtime/changelog.Debian.gz
/usr/share/doc/gnustep-base-runtime/changelog.gz
/usr/share/doc/gnustep-base-runtime/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/gnustep-base-runtime
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/HTMLLinker.1.gz
/usr/share/man/man1/autogsdoc.1.gz
/usr/share/man/man1/cvtenc.1.gz
/usr/share/man/man1/defaults.1.gz
/usr/share/man/man1/gdnc.1.gz
/usr/share/man/man1/gspath.1.gz
/usr/share/man/man1/make_strings.1.gz
/usr/share/man/man1/pldes.1.gz
/usr/share/man/man1/plutil.1.gz
/usr/share/man/man1/sfparse.1.gz
/usr/share/man/man1/xmlparse.1.gz
/usr/share/man/man8
/usr/share/man/man8/gdomap.8.gz
/usr/share/man/man1/pl.1.gz
/usr/share/man/man1/pl2link.1.gz
/usr/share/man/man1/plget.1.gz
/usr/share/man/man1/plmerge.1.gz
/usr/share/man/man1/plparse.1.gz
/usr/share/man/man1/plser.1.gz
