# src/var/lib/dpkg/info/live-tools.list
/.
/bin
/bin/live-medium-cache
/bin/live-medium-eject
/bin/live-partial-squashfs-updates
/bin/live-persistence
/bin/live-system
/bin/live-toram
/bin/live-update-initramfs
/bin/live-update-initramfs-uuid
/etc
/etc/init.d
/etc/init.d/live-tools
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/live-tools.service
/lib/systemd/system-shutdown
/lib/systemd/system-shutdown/live-tools.shutdown
/usr
/usr/sbin
/usr/share
/usr/share/bug
/usr/share/bug/live-tools
/usr/share/bug/live-tools/presubj
/usr/share/doc
/usr/share/doc/live-tools
/usr/share/doc/live-tools/changelog.gz
/usr/share/doc/live-tools/copyright
/usr/share/man
/usr/share/man/ca
/usr/share/man/ca/man7
/usr/share/man/ca/man7/live-tools.7.gz
/usr/share/man/ca/man8
/usr/share/man/ca/man8/live-update-initramfs.8.gz
/usr/share/man/de
/usr/share/man/de/man1
/usr/share/man/de/man7
/usr/share/man/de/man7/live-tools.7.gz
/usr/share/man/de/man8
/usr/share/man/de/man8/live-update-initramfs.8.gz
/usr/share/man/es
/usr/share/man/es/man7
/usr/share/man/es/man7/live-tools.7.gz
/usr/share/man/es/man8
/usr/share/man/es/man8/live-update-initramfs.8.gz
/usr/share/man/fr
/usr/share/man/fr/man7
/usr/share/man/fr/man7/live-tools.7.gz
/usr/share/man/fr/man8
/usr/share/man/fr/man8/live-update-initramfs.8.gz
/usr/share/man/ja
/usr/share/man/ja/man7
/usr/share/man/ja/man7/live-tools.7.gz
/usr/share/man/ja/man8
/usr/share/man/ja/man8/live-update-initramfs.8.gz
/usr/share/man/man1
/usr/share/man/man7
/usr/share/man/man7/live-tools.7.gz
/usr/share/man/man8
/usr/share/man/man8/live-update-initramfs.8.gz
/usr/sbin/update-initramfs
/usr/share/man/de/man1/live-system.1.gz
/usr/share/man/de/man1/live-toram.1.gz
/usr/share/man/man1/live-system.1.gz
/usr/share/man/man1/live-toram.1.gz
/usr/share/man/man8/update-initramfs.8.gz
