# src/var/lib/dpkg/info/memlockd.list
/.
/etc
/etc/default
/etc/default/memlockd
/etc/init.d
/etc/init.d/memlockd
/etc/memlockd.cfg
/etc/memlockd.d
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/memlockd.service
/usr
/usr/sbin
/usr/sbin/memlockd
/usr/share
/usr/share/doc
/usr/share/doc/memlockd
/usr/share/doc/memlockd/changelog.Debian.amd64.gz
/usr/share/doc/memlockd/changelog.Debian.gz
/usr/share/doc/memlockd/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/memlockd.8.gz
