# src/var/lib/dpkg/info/memlockd.conffiles
/etc/default/memlockd
/etc/init.d/memlockd
/etc/memlockd.cfg
