# src/var/lib/dpkg/info/plymouth.list
/.
/etc
/etc/init.d
/etc/init.d/plymouth
/etc/init.d/plymouth-log
/etc/logrotate.d
/etc/logrotate.d/bootlog
/etc/plymouth
/etc/plymouth/plymouthd.conf
/lib
/lib/lsb
/lib/lsb/init-functions.d
/lib/lsb/init-functions.d/99-plymouth
/lib/systemd
/lib/systemd/system
/lib/systemd/system/halt.target.wants
/lib/systemd/system/initrd-switch-root.target.wants
/lib/systemd/system/kexec.target.wants
/lib/systemd/system/multi-user.target.wants
/lib/systemd/system/plymouth-halt.service
/lib/systemd/system/plymouth-kexec.service
/lib/systemd/system/plymouth-poweroff.service
/lib/systemd/system/plymouth-quit-wait.service
/lib/systemd/system/plymouth-quit.service
/lib/systemd/system/plymouth-read-write.service
/lib/systemd/system/plymouth-reboot.service
/lib/systemd/system/plymouth-start.service
/lib/systemd/system/plymouth-switch-root-initramfs.service
/lib/systemd/system/plymouth-switch-root.service
/lib/systemd/system/poweroff.target.wants
/lib/systemd/system/reboot.target.wants
/lib/systemd/system/sysinit.target.wants
/lib/systemd/system/systemd-ask-password-plymouth.path
/lib/systemd/system/systemd-ask-password-plymouth.service
/usr
/usr/bin
/usr/bin/plymouth
/usr/lib
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/plymouth
/usr/lib/x86_64-linux-gnu/plymouth/details.so
/usr/lib/x86_64-linux-gnu/plymouth/renderers
/usr/lib/x86_64-linux-gnu/plymouth/renderers/drm.so
/usr/lib/x86_64-linux-gnu/plymouth/renderers/frame-buffer.so
/usr/lib/x86_64-linux-gnu/plymouth/script.so
/usr/lib/x86_64-linux-gnu/plymouth/text.so
/usr/lib/x86_64-linux-gnu/plymouth/tribar.so
/usr/libexec
/usr/libexec/plymouth
/usr/libexec/plymouth/plymouth-generate-initrd
/usr/libexec/plymouth/plymouth-populate-initrd
/usr/libexec/plymouth/plymouth-update-initrd
/usr/libexec/plymouth/plymouthd-fd-escrow
/usr/sbin
/usr/sbin/plymouth-set-default-theme
/usr/sbin/plymouthd
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_plymouth.py
/usr/share/doc
/usr/share/doc/plymouth
/usr/share/doc/plymouth/NEWS.Debian.gz
/usr/share/doc/plymouth/README
/usr/share/doc/plymouth/README.Debian
/usr/share/doc/plymouth/TODO
/usr/share/doc/plymouth/TODO.Debian
/usr/share/doc/plymouth/changelog.Debian.gz
/usr/share/doc/plymouth/changelog.gz
/usr/share/doc/plymouth/copyright
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/plymouth
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/plymouth
/usr/share/initramfs-tools/scripts/init-premount
/usr/share/initramfs-tools/scripts/init-premount/plymouth
/usr/share/initramfs-tools/scripts/panic
/usr/share/initramfs-tools/scripts/panic/plymouth
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/plymouth
/usr/share/locale
/usr/share/locale/be
/usr/share/locale/be/LC_MESSAGES
/usr/share/locale/be/LC_MESSAGES/plymouth.mo
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/cs/LC_MESSAGES/plymouth.mo
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/da/LC_MESSAGES/plymouth.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/plymouth.mo
/usr/share/locale/eo
/usr/share/locale/eo/LC_MESSAGES
/usr/share/locale/eo/LC_MESSAGES/plymouth.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/plymouth.mo
/usr/share/locale/eu
/usr/share/locale/eu/LC_MESSAGES
/usr/share/locale/eu/LC_MESSAGES/plymouth.mo
/usr/share/locale/fa
/usr/share/locale/fa/LC_MESSAGES
/usr/share/locale/fa/LC_MESSAGES/plymouth.mo
/usr/share/locale/fi
/usr/share/locale/fi/LC_MESSAGES
/usr/share/locale/fi/LC_MESSAGES/plymouth.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/plymouth.mo
/usr/share/locale/fur
/usr/share/locale/fur/LC_MESSAGES
/usr/share/locale/fur/LC_MESSAGES/plymouth.mo
/usr/share/locale/he
/usr/share/locale/he/LC_MESSAGES
/usr/share/locale/he/LC_MESSAGES/plymouth.mo
/usr/share/locale/hr
/usr/share/locale/hr/LC_MESSAGES
/usr/share/locale/hr/LC_MESSAGES/plymouth.mo
/usr/share/locale/hu
/usr/share/locale/hu/LC_MESSAGES
/usr/share/locale/hu/LC_MESSAGES/plymouth.mo
/usr/share/locale/id
/usr/share/locale/id/LC_MESSAGES
/usr/share/locale/id/LC_MESSAGES/plymouth.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/plymouth.mo
/usr/share/locale/lt
/usr/share/locale/lt/LC_MESSAGES
/usr/share/locale/lt/LC_MESSAGES/plymouth.mo
/usr/share/locale/ms
/usr/share/locale/ms/LC_MESSAGES
/usr/share/locale/ms/LC_MESSAGES/plymouth.mo
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/nl/LC_MESSAGES/plymouth.mo
/usr/share/locale/pa
/usr/share/locale/pa/LC_MESSAGES
/usr/share/locale/pa/LC_MESSAGES/plymouth.mo
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/plymouth.mo
/usr/share/locale/pt
/usr/share/locale/pt/LC_MESSAGES
/usr/share/locale/pt/LC_MESSAGES/plymouth.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/plymouth.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/plymouth.mo
/usr/share/locale/si
/usr/share/locale/si/LC_MESSAGES
/usr/share/locale/si/LC_MESSAGES/plymouth.mo
/usr/share/locale/sk
/usr/share/locale/sk/LC_MESSAGES
/usr/share/locale/sk/LC_MESSAGES/plymouth.mo
/usr/share/locale/sr
/usr/share/locale/sr/LC_MESSAGES
/usr/share/locale/sr/LC_MESSAGES/plymouth.mo
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/plymouth.mo
/usr/share/locale/tr
/usr/share/locale/tr/LC_MESSAGES
/usr/share/locale/tr/LC_MESSAGES/plymouth.mo
/usr/share/locale/uk
/usr/share/locale/uk/LC_MESSAGES
/usr/share/locale/uk/LC_MESSAGES/plymouth.mo
/usr/share/locale/zh_CN
/usr/share/locale/zh_CN/LC_MESSAGES
/usr/share/locale/zh_CN/LC_MESSAGES/plymouth.mo
/usr/share/locale/zh_HK
/usr/share/locale/zh_HK/LC_MESSAGES
/usr/share/locale/zh_HK/LC_MESSAGES/plymouth.mo
/usr/share/locale/zh_TW
/usr/share/locale/zh_TW/LC_MESSAGES
/usr/share/locale/zh_TW/LC_MESSAGES/plymouth.mo
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/plymouth-set-default-theme.1.gz
/usr/share/man/man1/plymouth.1.gz
/usr/share/man/man8
/usr/share/man/man8/plymouth.8.gz
/usr/share/man/man8/plymouthd.8.gz
/usr/share/plymouth
/usr/share/plymouth/debian-logo.png
/usr/share/plymouth/plymouthd.defaults
/usr/share/plymouth/themes
/usr/share/plymouth/themes/details
/usr/share/plymouth/themes/details/details.plymouth
/usr/share/plymouth/themes/text
/usr/share/plymouth/themes/text/text.plymouth
/usr/share/plymouth/themes/tribar
/usr/share/plymouth/themes/tribar/tribar.plymouth
/var
/var/lib
/var/lib/plymouth
/lib/systemd/system/halt.target.wants/plymouth-halt.service
/lib/systemd/system/halt.target.wants/plymouth-switch-root-initramfs.service
/lib/systemd/system/initrd-switch-root.target.wants/plymouth-start.service
/lib/systemd/system/initrd-switch-root.target.wants/plymouth-switch-root.service
/lib/systemd/system/kexec.target.wants/plymouth-kexec.service
/lib/systemd/system/kexec.target.wants/plymouth-switch-root-initramfs.service
/lib/systemd/system/multi-user.target.wants/plymouth-quit-wait.service
/lib/systemd/system/multi-user.target.wants/plymouth-quit.service
/lib/systemd/system/plymouth-log.service
/lib/systemd/system/plymouth.service
/lib/systemd/system/poweroff.target.wants/plymouth-poweroff.service
/lib/systemd/system/poweroff.target.wants/plymouth-switch-root-initramfs.service
/lib/systemd/system/reboot.target.wants/plymouth-reboot.service
/lib/systemd/system/reboot.target.wants/plymouth-switch-root-initramfs.service
/lib/systemd/system/sysinit.target.wants/plymouth-read-write.service
/lib/systemd/system/sysinit.target.wants/plymouth-start.service
/usr/share/man/man8/plymouth-set-default-theme.8.gz
