# src/var/lib/dpkg/info/rng-tools-debian.list
/.
/etc
/etc/default
/etc/default/rng-tools-debian
/etc/init.d
/etc/init.d/rng-tools-debian
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/rng-tools-debian
/etc/logcheck/violations.ignore.d
/etc/logcheck/violations.ignore.d/rng-tools-debian
/usr
/usr/bin
/usr/bin/rngtest
/usr/sbin
/usr/sbin/rngd
/usr/share
/usr/share/doc
/usr/share/doc/rng-tools-debian
/usr/share/doc/rng-tools-debian/NEWS.Debian.gz
/usr/share/doc/rng-tools-debian/NEWS.gz
/usr/share/doc/rng-tools-debian/README.gz
/usr/share/doc/rng-tools-debian/changelog.gz
/usr/share/doc/rng-tools-debian/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/rng-tools-debian
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/rngtest.1.gz
/usr/share/man/man8
/usr/share/man/man8/rngd.8.gz
