# src/var/lib/dpkg/info/systemd-oomd.list
/.
/etc
/etc/systemd
/etc/systemd/oomd.conf
/usr
/usr/bin
/usr/bin/oomctl
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/oomd.conf.d
/usr/lib/systemd/oomd.conf.d/10-oomd-defaults.conf
/usr/lib/systemd/system
/usr/lib/systemd/system/-.slice.d
/usr/lib/systemd/system/-.slice.d/10-oomd-root-slice-defaults.conf
/usr/lib/systemd/system/systemd-oomd.service
/usr/lib/systemd/system/systemd-oomd.socket
/usr/lib/systemd/system/user@.service.d
/usr/lib/systemd/system/user@.service.d/10-oomd-user-service-defaults.conf
/usr/lib/systemd/systemd-oomd
/usr/lib/sysusers.d
/usr/lib/sysusers.d/systemd-oom.conf
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/oomctl
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.oom1.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/org.freedesktop.oom1.conf
/usr/share/doc
/usr/share/doc/systemd-oomd
/usr/share/doc/systemd-oomd/NEWS.Debian.gz
/usr/share/doc/systemd-oomd/changelog.Debian.gz
/usr/share/doc/systemd-oomd/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/systemd-oomd
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/oomctl.1.gz
/usr/share/man/man5
/usr/share/man/man5/oomd.conf.5.gz
/usr/share/man/man5/org.freedesktop.oom1.5.gz
/usr/share/man/man8
/usr/share/man/man8/systemd-oomd.service.8.gz
/usr/share/zsh
/usr/share/zsh/vendor-completions
/usr/share/zsh/vendor-completions/_oomctl
/usr/share/man/man5/oomd.conf.d.5.gz
/usr/share/man/man8/systemd-oomd.8.gz
