# src/var/lib/dpkg/info/tor.list
/.
/etc
/etc/apparmor.d
/etc/apparmor.d/abstractions
/etc/apparmor.d/abstractions/tor
/etc/apparmor.d/system_tor
/etc/cron.weekly
/etc/cron.weekly/tor
/etc/default
/etc/default/tor
/etc/init.d
/etc/init.d/tor
/etc/logrotate.d
/etc/logrotate.d/tor
/etc/runit
/etc/runit/runsvdir
/etc/runit/runsvdir/default
/etc/sv
/etc/sv/tor
/etc/sv/tor/.meta
/etc/sv/tor/.meta/installed
/etc/sv/tor/log
/etc/sv/tor/log/run
/etc/sv/tor/run
/etc/tor
/etc/tor/torrc
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/tor.service
/lib/systemd/system/tor@.service
/lib/systemd/system/tor@default.service
/lib/systemd/system-generators
/lib/systemd/system-generators/tor-generator
/usr
/usr/bin
/usr/bin/tor
/usr/bin/tor-gencert
/usr/bin/tor-print-ed-signing-cert
/usr/bin/tor-resolve
/usr/bin/torify
/usr/sbin
/usr/sbin/tor-instance-create
/usr/share
/usr/share/doc
/usr/share/doc/tor
/usr/share/doc/tor/NEWS.Debian.gz
/usr/share/doc/tor/README.Debian
/usr/share/doc/tor/changelog.Debian.gz
/usr/share/doc/tor/changelog.gz
/usr/share/doc/tor/copyright
/usr/share/doc/tor/tor-exit-notice.html
/usr/share/doc/tor/tor-gencert.html
/usr/share/doc/tor/tor-print-ed-signing-cert.html
/usr/share/doc/tor/tor-resolve.html
/usr/share/doc/tor/tor.html
/usr/share/doc/tor/torify.html
/usr/share/doc/tor/torrc.sample.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/tor
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/tor-gencert.1.gz
/usr/share/man/man1/tor-print-ed-signing-cert.1.gz
/usr/share/man/man1/tor-resolve.1.gz
/usr/share/man/man1/tor.1.gz
/usr/share/man/man1/torify.1.gz
/usr/share/man/man5
/usr/share/man/man8
/usr/share/man/man8/tor-instance-create.8.gz
/usr/share/runit
/usr/share/runit/meta
/usr/share/runit/meta/tor
/usr/share/runit/meta/tor/installed
/usr/share/tor
/usr/share/tor/tor-service-defaults-torrc
/usr/share/tor/tor-service-defaults-torrc-instances
/var
/var/log
/var/log/runit
/var/log/runit/tor
/usr/sbin/tor
/usr/share/man/man5/torrc.5.gz
