# src/var/lib/dpkg/info/usbguard.conffiles
/etc/init.d/usbguard
/etc/init.d/usbguard-dbus
/etc/logrotate.d/usbguard
/etc/usbguard/usbguard-daemon.conf
