# src/var/lib/dpkg/info/usbguard.list
/.
/etc
/etc/init.d
/etc/init.d/usbguard
/etc/init.d/usbguard-dbus
/etc/logrotate.d
/etc/logrotate.d/usbguard
/etc/usbguard
/etc/usbguard/IPCAccessControl.d
/etc/usbguard/rules.d
/etc/usbguard/usbguard-daemon.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/usbguard-dbus.service
/lib/systemd/system/usbguard.service
/usr
/usr/bin
/usr/bin/usbguard
/usr/bin/usbguard-rule-parser
/usr/sbin
/usr/sbin/usbguard-daemon
/usr/sbin/usbguard-dbus
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/usbguard
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.usbguard1.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/org.usbguard1.conf
/usr/share/doc
/usr/share/doc/usbguard
/usr/share/doc/usbguard/NEWS.Debian.gz
/usr/share/doc/usbguard/README.Debian
/usr/share/doc/usbguard/changelog.Debian.amd64.gz
/usr/share/doc/usbguard/changelog.Debian.gz
/usr/share/doc/usbguard/changelog.gz
/usr/share/doc/usbguard/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/usbguard
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/usbguard.1.gz
/usr/share/man/man5
/usr/share/man/man5/usbguard-daemon.conf.5.gz
/usr/share/man/man5/usbguard-rules.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/usbguard-daemon.8.gz
/usr/share/man/man8/usbguard-dbus.8.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.usbguard1.policy
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/org.usbguard1.rules
/usr/share/zsh
/usr/share/zsh/vendor-completions
/usr/share/zsh/vendor-completions/_usbguard
/var
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/org.usbguard1.pkla
/var/log
/var/log/usbguard
