# src/var/lib/dpkg/info/xdg-desktop-portal.postrm
#!/bin/sh
set -e
# Automatically added by dh_installsystemduser/13.11.4
if [ "$1" = "purge" ]; then
	if [ -z "${DPKG_ROOT:-}" ] && [ -x "/usr/bin/deb-systemd-helper" ] ; then
		deb-systemd-helper --user purge 'xdg-desktop-portal-rewrite-launchers.service' >/dev/null || true
	fi
fi
# End automatically added section
