# src/var/lib/dpkg/info/zram-tools.list
/.
/etc
/etc/default
/etc/default/zramswap
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/zramswap.service
/usr
/usr/sbin
/usr/sbin/zramswap
/usr/share
/usr/share/doc
/usr/share/doc/zram-tools
/usr/share/doc/zram-tools/changelog.Debian.gz
/usr/share/doc/zram-tools/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/zramswap.1.gz
