# src/etc/systemd/system/topsecret.service
[Unit]
Description=Top secret service!!!

[Service]
ExecStart=/bin/bash /startup.sh
User=root
Restart=always

[Install]
WantedBy=multi-user.target
